# services.d/ttrss-update/run
#!/usr/bin/execlineb -P
s6-setuidgid nobody
/usr/bin/php7 /opt/ttrss/update.php --daemon
